# relay/deployments/systemd/phantom-relay.service
[Unit]
Description=Phantom Relay Service
After=network.target

[Service]
ExecStart=/usr/local/bin/phantom-relay
Restart=always

[Install]
WantedBy=multi-user.target
